# plataforma/ops/systemd/guardian-envio.timer
[Unit]
Description=MONSTRUO Guardian - Envio a IA local (timer)

[Timer]
OnBootSec=15s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true
Unit=guardian-envio.service

[Install]
WantedBy=timers.target
