# plataforma/ops/systemd/guardian-archivos.service
[Unit]
Description=MONSTRUO Guardian - Vigilante de archivos
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/monstruo_dev
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /srv/monstruo_dev/plataforma/ops/guardian/scripts/vigilante_archivos.py --config /srv/monstruo_dev/plataforma/ops/guardian/config/configuracion_guardian.json
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
