# ops/systemd/guardian-limpieza.timer
[Unit]
Description=MONSTRUO Guardian - Limpieza de backups (timer)

[Timer]
OnCalendar=*-*-* 03:25:00
RandomizedDelaySec=30s
Persistent=true
Unit=monstruo-guardian-limpieza-backups.service

[Install]
WantedBy=timers.target
