# ops/systemd/guardian-limpieza.service
[Unit]
Description=MONSTRUO Guardian - Limpieza de backups (oneshot)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/monstruo
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/monstruo/venv/bin/python /srv/monstruo/ops/guardian/limpiador_backups.py --config /srv/monstruo/ops/guardian/configuracion_guardian.json --once
