# ops/systemd/guardian-registros.service
[Unit]
Description=MONSTRUO Guardian - Vigilante de registros
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/monstruo
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/monstruo/venv/bin/python /srv/monstruo/ops/guardian/vigilante_registros.py --config /srv/monstruo/ops/guardian/configuracion_guardian.json
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
