# ops/systemd/guardian-envio.service
[Unit]
Description=MONSTRUO Guardian - Envio a IA local (oneshot)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/monstruo_dev
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/srv/monstruo_dev/ops/guardian/config/entorno_guardian.env
ExecStart=/usr/bin/python3 /srv/monstruo_dev/ops/guardian/scripts/enviar_a_ia_local.py --config /srv/monstruo_dev/ops/guardian/config/configuracion_guardian.json --once --limite 10 --timeout 3
